# influxdb/rootfs/etc/services.d/bootstrap/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: InfluxDB 3
# Keep a placeholder bootstrap service without probing InfluxDB during startup
# ==============================================================================
set -euo pipefail

bashio::log.info "Skipping automatic database bootstrap during startup."

exec sleep infinity
