# influxdb/rootfs/etc/services.d/bootstrap/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: InfluxDB 3
# Bootstrap default database
# ==============================================================================
set -euo pipefail

export LD_LIBRARY_PATH="/usr/lib/influxdb3/python/lib:${LD_LIBRARY_PATH:-}"
export PYTHONHOME="/usr/lib/influxdb3/python"
export PYTHONPATH="/usr/lib/influxdb3/python/lib/python3.13"
declare token
declare database
declare host_url
declare list_url
declare create_url

s6-svwait -u -t 5000 /run/service//influxdb
host_url="http://127.0.0.1:8181"
list_url="${host_url}/api/v3/configure/database?format=json"
create_url="${host_url}/api/v3/configure/database"

database="$(bashio::config 'default_database')"
if ! bashio::var.has_value "${database}"; then
    database="homeassistant"
fi

if bashio::fs.file_exists "/data/influxdb3/admin-token.txt"; then
    token="$(cat /data/influxdb3/admin-token.txt)"
else
    bashio::log.warning "Admin token not found; skipping database bootstrap."
    exit 0
fi

# Wait for the HTTP listener before issuing CLI requests to the server.
if ! bashio::net.wait_for 8181 localhost 120000; then
    bashio::log.warning "InfluxDB 3 HTTP API did not become reachable, skipping database bootstrap."
    exit 0
fi

for i in {120..0}; do
    if curl --silent --show-error --fail \
        --header "Authorization: Bearer ${token}" \
        "${list_url}" \
        > /tmp/influxdb3-dbs.json \
        2>/tmp/influxdb3-dbs.err; then
        break
    fi
    if (( i % 5 == 0 )); then
        if [[ -s /tmp/influxdb3-dbs.err ]]; then
            bashio::log.info "Waiting for InfluxDB 3 to become ready... ($(head -n 1 /tmp/influxdb3-dbs.err))"
        else
            bashio::log.info "Waiting for InfluxDB 3 to become ready..."
        fi
    fi
    sleep 2
done

if [[ "${i}" = 0 ]]; then
    bashio::log.warning "InfluxDB 3 did not become ready, skipping database bootstrap."
    exit 0
fi

if grep -q "\"iox::database\"[[:space:]]*:[[:space:]]*\"${database}\"" /tmp/influxdb3-dbs.json; then
    bashio::log.info "Database '${database}' already exists."
else
    bashio::log.info "Creating database '${database}'..."
    curl --silent --show-error --fail \
        --request POST \
        --header "Authorization: Bearer ${token}" \
        --header "Content-Type: application/json" \
        --data "{\"db\":\"${database}\"}" \
        "${create_url}" \
        > /dev/null 2>&1 || \
        bashio::log.warning "Failed to create database '${database}'."
fi

exec sleep infinity
