# influxdb/rootfs/etc/services.d/influxdb/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: InfluxDB 3
# Runs the InfluxDB 3 server
# ==============================================================================
set -euo pipefail

declare edition
declare cluster_id
declare node_id
declare license_email
declare license_file
declare license_type
declare token_file
declare binary
declare log_filter
declare -a args

export LD_LIBRARY_PATH="/usr/lib/influxdb3/python/lib:${LD_LIBRARY_PATH:-}"

edition="$(bashio::config 'edition')"
license_email="$(bashio::config 'license_email')"
license_file="$(bashio::config 'license_file')"
license_type="$(bashio::config 'license_type')"
cluster_id="$(bashio::config 'cluster_id')"
node_id="$(bashio::config 'node_id')"
token_file="/data/influxdb3/admin-token.json"

if [[ -n "${license_email}" ]] || [[ -n "${license_file}" ]]; then
    edition="enterprise"
fi
if [[ -z "${edition}" ]]; then
    edition="core"
fi

binary="/usr/local/bin/influxdb3-${edition}"
if ! bashio::fs.file_exists "${binary}"; then
    bashio::exit.nok "InfluxDB 3 ${edition} binary missing at ${binary}"
fi

if [[ "${edition}" = "enterprise" ]] && [[ -z "${license_email}" ]] && [[ -z "${license_file}" ]]; then
    bashio::exit.nok "Enterprise edition selected but no license_email or license_file set."
fi

args=(serve)
args+=(--node-id "${node_id:-ha-node}")
args+=(--object-store file)
args+=(--data-dir /data/influxdb3/data)
args+=(--plugin-dir /data/influxdb3/plugins)

if bashio::config.false 'auth'; then
    bashio::log.warning "InfluxDB authentication is disabled!"
    bashio::log.warning "This is NOT recommended!!!"
    args+=(--without-auth)
else
    args+=(--admin-token-file "${token_file}")
fi

if bashio::config.has_value 'log_level'; then
    case "$(bashio::string.lower "$(bashio::config 'log_level')")" in
        trace)
            log_filter="trace"
            ;;
        debug)
            log_filter="debug"
            ;;
        info|notice)
            log_filter="info"
            ;;
        warning)
            log_filter="warn"
            ;;
        error|fatal)
            log_filter="error"
            ;;
    esac
    if bashio::var.has_value "${log_filter}"; then
        args+=(--log-filter "${log_filter}")
    fi
fi

if [ "${edition}" = "enterprise" ]; then
    args+=(--cluster-id "${cluster_id:-ha-cluster}")
    if bashio::var.has_value "${license_file}"; then
        args+=(--license-file "${license_file}")
    elif bashio::var.has_value "${license_email}"; then
        args+=(--license-email "${license_email}")
        args+=(--license-type "${license_type:-home}")
    fi
fi

bashio::log.info "Starting InfluxDB 3 (${edition})..."
exec "${binary}" "${args[@]}"
